# config/autostart/gnome-hide-panel.desktop
[Desktop Entry]
Type=Application
Name=GNOME Hide Panel
Exec=/home/user/scripts/gnome-hide-panel.sh
