# config/autostart/nvidia-powermizer.desktop
[Desktop Entry]
Type=Application
Name=nvidia-powermizer
Exec=nvidia-settings -a GpuPowerMizerMode=1

